# showcase.pre.list
./Sync/rne/mount_karotaxi.sh
./Sync/rne/.local/share/Trash/files/usefulctl.sh
./Pobrane/quake2-7.01/stuff/cdripper.sh
./Pobrane/quake2-7.01/stuff/quake2-start.sh
./keystate-half-fix.sh
./homefiles/.myrc.d/namespaces.sh
./homefiles/.myrc.d/rcops.sh
./homefiles/.myrc.d/choose.sh
./homefiles/.myrc.d/disabled/12_fileops.sh
./homefiles/.myrc.d/disabled/11_pkgmanage.sh
./homefiles/.myrc.d/disabled/21_usefulpkg.sh
./homefiles/.myrc.d/disabled/Z_test.sh
./homefiles/.myrc.d/disabled/tui.sh
./homefiles/.myrc.d/.backup/namespaces.sh
./homefiles/.myrc.d/.backup/rcops.sh
./homefiles/.myrc.d/.backup/choose.sh
./homefiles/.myrc.d/.backup/misc.sh
./homefiles/.myrc.d/.backup/math.sh
./homefiles/.myrc.d/.backup/wrappers.sh
./homefiles/.myrc.d/.backup/markfile.sh
./homefiles/.myrc.d/.backup/review.sh
./homefiles/.myrc.d/.backup/tui.sh
./homefiles/.myrc.d/.backup/logbook.sh
./homefiles/.myrc.d/.backup/output.sh
./homefiles/.myrc.d/.backup/bashtools.sh
./homefiles/.myrc.d/.backup/permavars.sh
./homefiles/.myrc.d/misc.sh
./homefiles/.myrc.d/math.sh
./homefiles/.myrc.d/wrappers.sh
./homefiles/.myrc.d/deprecated/12_archives.sh
./homefiles/.myrc.d/deprecated/install.sh
./homefiles/.myrc.d/markfile.sh
./homefiles/.myrc.d/review.sh
./homefiles/.myrc.d/tests/bash.sh
./homefiles/.myrc.d/tests/file:tag.sh
./homefiles/.myrc.d/tests/logbook.sh
./homefiles/.myrc.d/tests/bash.sh.out
./homefiles/.myrc.d/logbook.sh
./homefiles/.myrc.d/output.sh
./homefiles/.myrc.d/bashtools.sh
./homefiles/.myrc.d/permavars.sh
./homefiles/src/myrc/modules/file/tag/meta.sh
./homefiles/projectyellowcab/mount.sh
./homefiles/.iep/dummyFiles/dummy.sh
./homefiles/useful/global/r.sh
./homefiles/useful/libs/whiptailctx.sh
./homefiles/useful/scripts/jpgcomment.sh
./homefiles/useful/scripts/findnew.sh
./homefiles/useful/scripts/valuables.sh
./homefiles/useful/scripts/toggle_mouse.sh
./homefiles/useful/scripts/count-minutes.sh
./homefiles/useful/scripts/shsession.sh
./alicia.sh
./multimachinehead.sh
./remoterun.sh
./src/Bequnge/mkbundle.sh
./src/xleapmouse/config.sh
./src/xleapmouse/ltmain.sh
./tablet-init.sh
./byond.sh
./mnt/bukkit/runserver.sh
./mnt/bukkit/eu_spigot/runserver.sh
./.my/enable.sh
./programming/projekty/maritime/crawl/crawlcontinue.sh
./programming/projekty/maritime/crawl/crawler/rawtoclean.sh
./programming/projekty/maritime/crawl/crawler/crawl.sh
./programming/projekty/maritime/crawl/crawl.sh
./programming/projekty/maritime/crawl/fieldneigh.sh
./programming/projekty/watcher_test/tac.sh
./programming/povray-samples/camera/mesh_camera/bake.sh
./programming/pycharm-community-2016.3.2/bin/inspect.sh
./programming/pycharm-community-2016.3.2/bin/pycharm.sh
./programming/pycharm-community-2016.3.2/bin/format.sh
./programming/pycharm-community-2016.3.2/helpers/pydev/pydevd_attach_to_process/linux/compile_so.sh
./programming/pycharm-community-2016.3.2/helpers/pydev/pydevd_attach_to_process/linux/compile_mac.sh
./programming/befunge/fuzzydistance/watcher.sh
./programming/8th/bin/floatsign.sh
./keystate.sh
